# systemd/leaderboard.timer
[Unit]
Description=Update the Advent of Code leaderboard

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
